# dev-db/influxdb-bin/files/influxd.rc
#!/sbin/openrc-run

name="InfluxDB"
pidfile="/run/influxd.pid"
command="/usr/bin/influxd"
command_background="yes"
command_user="influxdb"
